# server/etc/cron.weekly/rs-backup-rotate
#!/bin/bash
/usr/sbin/rs-rotate-cron weekly
