# server/etc/cron.monthly/rs-backup-rotate
#!/bin/bash
/usr/sbin/rs-rotate-cron monthly
